# pip-bot.service
[Unit]
Description=pip-bot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pip-bot
ExecStart=/home/pi/.local/bin/poetry run python -m bot
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
# Resource limits to protect Raspberry Pi from OOM and runaway processes
MemoryLimit=800M
OOMPolicy=kill
TasksMax=50
CPUQuota=80%

[Install]
WantedBy=multi-user.target
